# game-v2/systemd/game-v2.service
[Unit]
Description=Clash Demo v2 (Fastify + WS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/srv/game-v2/server
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8080
ExecStart=/usr/bin/node /srv/game-v2/server/dist/index.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
